# flexsearch/mcp_server.service
# Flexsearch MCP Server - systemd user service
#
# Install:
#   cp flexsearch/mcp_server.service ~/.config/systemd/user/flexsearch-mcp.service
#   systemctl --user daemon-reload
#   systemctl --user enable flexsearch-mcp
#   systemctl --user start flexsearch-mcp
#
# Then expose via Tailscale Funnel:
#   tailscale funnel 8080
#
# View logs:
#   journalctl --user -u flexsearch-mcp -f

[Unit]
Description=Flexsearch MCP Server
After=network.target

[Service]
Type=simple
ExecStart=/home/axp/projects/flexsearch/main/venv/bin/python -m flexsearch.mcp_server --http --port 8081
WorkingDirectory=/home/axp/projects/flexsearch/main
Environment="PYTHONPATH=/home/axp/projects/flexsearch/main"

Restart=always
RestartSec=5

[Install]
WantedBy=default.target
